# data/systemd/systemd-coredump@.service.d/apport-coredump-hook.conf
[Unit]
OnSuccess=apport-coredump-hook@%i.service
