# debian/apport.cron.daily
#!/bin/sh -e
# clean all crash reports which are older than a week.
find /var/crash -mtime +7 -print0 | xargs -0 rm 
