# etc/init/iwlwifi-error-dump.conf
start on ieee80211-device-changed DRIVER=iwlwifi EVENT=error_dump

script
    phyname=$(basename ${DEVPATH})
    /usr/share/apport/iwlwifi_error_dump "$phyname"
end script
